# code/my-notifier/my-notifier.service
[Unit]
Description=My Notifier

[Service]
Type=notify
ExecStart=/usr/local/bin/my-notifier.sh
