# code/my-watchdog/my-watchdog.service
[Unit]
Description=My Watchdog

[Service]
Type=notify
WatchdogSec=10s
ExecStart=/usr/local/bin/test-watchdog
