# code/my-memory-monitor/my-memory-monitor.service
[Unit]
Description=My Memory Monitor

[Service]
ExecStart=/usr/local/bin/my-memory-monitor.sh
