# code/my-logger/my-logger.service
[Unit]
Description=My Logger

[Service]
ExecStart=/usr/local/bin/my-logger.sh
