# code/my-logger/my-logger.socket
[Unit]
Description=Socket for My Logger

[Socket]
ListenFIFO=/tmp/my-logger.fifo

[Install]
WantedBy=sockets.target
